# resources/ec2-termination-handler.service
[Unit]
Description=EC2 termination handler

[Service]
ExecStart=/usr/local/bin/ec2-termination-handler-rust

[Install]
WantedBy=multi-user.target
